# systemd/adminpanelaz-node.service
[Unit]
Description=AdminPanelAZ Node Agent (AntiZapret remote node)
Documentation=file:///opt/AdminPanelAZ/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/AdminPanelAZ
EnvironmentFile=-/opt/AdminPanelAZ/backend/node_agent.env
Environment=NODE_AGENT_MODE=prod
Environment=NODE_AGENT_STATE_DIR=/var/lib/adminpanelaz-node
Environment=NODE_AGENT_HOST=0.0.0.0
Environment=NODE_AGENT_PORT=9100
Environment=NODE_AGENT_API_KEY=change-me-node-agent-key
Environment=ANTIZAPRET_PATH=/root/antizapret
ExecStart=/opt/AdminPanelAZ/start_node_agent.sh watchdog prod
ExecStop=/opt/AdminPanelAZ/start_node_agent.sh stop
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=adminpanelaz-node

[Install]
WantedBy=multi-user.target
